# god.conf
description "god"

start on runlevel [2]
stop on runlevel [016]

console owner

exec /var/www/.rbenv/shims/god -c /var/www/bpan-org/unicorn.god

respawn